# deploy/gexwheel-mentions.container
# Quadlet unit - place in ~/.config/containers/systemd/
# Generates gexwheel-mentions.service (run-to-completion, fired by its timer).
[Unit]
Description=gexwheel - daily Reddit mention scan

[Container]
Image=localhost/gexwheel:latest
Exec=mentions
Volume=%h/gexwheel-data:/data:Z
Environment=GEXWHEEL_CONFIG=/data/config.yaml

[Service]
Type=oneshot
RemainAfterExit=no

[Install]
WantedBy=default.target
