# deploy/gexwheel-morning.container
# Quadlet unit - place in ~/.config/containers/systemd/
[Unit]
Description=gexwheel - weekday pre-market GEX screen + Discord alerts

[Container]
Image=localhost/gexwheel:latest
Exec=morning
Volume=%h/gexwheel-data:/data:Z
Environment=GEXWHEEL_CONFIG=/data/config.yaml

[Service]
Type=oneshot
RemainAfterExit=no

[Install]
WantedBy=default.target
