# deploy/gexwheel-morning.timer
# Place in ~/.config/systemd/user/
[Unit]
Description=Run gexwheel morning screen weekdays pre-market

[Timer]
OnCalendar=Mon..Fri 07:15:00 America/New_York
Persistent=true
Unit=gexwheel-morning.service

[Install]
WantedBy=timers.target
